# ops/live-studio-relay.service
[Unit]
Description=Live Studio RTMP Relay
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/live-studio
EnvironmentFile=/etc/live-studio/relay.env
ExecStart=/usr/bin/npm run relay
Restart=always
RestartSec=3
User=live-relay
Environment=NODE_ENV=production
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/live-studio /tmp
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
